# sup.conf
[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[inet_http_server]
port = 127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=http://127.0.0.1:9001

[program:anomaly]
command=python run.py configs/prod.prod.yml
process_name=%(program_name)s
numprocs=1
directory=/home/ubuntu/anomaly-detection/
autostart=true
startsecs=1
startretries=3
stdout_logfile=/var/log/anomaly.log

